# deploy/systemd/portfolio-backup.timer.example
[Unit]
Description=Run portfolio analytics backup daily

[Timer]
OnCalendar=*-*-* 00:30:00
Persistent=true
Unit=portfolio-backup.service

[Install]
WantedBy=timers.target
